# .cloud/php-fpm/supervisor.notification.ini

[program:workers_discovery]
command=bin/console app:pricing:worker discovery --no-debug -eprod
directory=/app
numprocs=10
autostart=true
autorestart=true
stderr_logfile=/var/log/app.workers.err.log
stdout_logfile=/var/log/app.workers.out.log
process_name=%(program_name)s_%(process_num)02d

[program:workers_pricing]
command=bin/console app:pricing:worker pricing --no-debug -eprod
directory=/app
numprocs=2
autostart=true
autorestart=true
stderr_logfile=/var/log/app.workers.err.log
stdout_logfile=/var/log/app.workers.out.log
process_name=%(program_name)s_%(process_num)02d
